# website.conf
description "website"

start on runlevel [2345]
stop on runlevel [016]

respawn
chdir /root/interactive-tutorials
exec gunicorn -b 0.0.0.0:80 -w 10 -t 30 main:app
